# unipi-modbus-tools/control/postinst
#!/bin/sh

systemctl --system daemon-reload
systemctl enable unipitcp.service
systemctl start unipitcp.service
