# unipi-modbus-tools/control/prerm
#!/bin/sh

systemctl stop unipitcp.service
systemctl disable unipitcp.service
systemctl --system daemon-reload
